# templates/udp2raw.service.j2
[Unit]
Description=udp2raw
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8060 -r 127.0.0.1:1250  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro & /opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8061 -r 127.0.0.1:1251  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro & /opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8062 -r 127.0.0.1:1252  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro & /opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8063 -r 127.0.0.1:1253  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro & /opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8064 -r 127.0.0.1:1254  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro & /opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8065 -r 127.0.0.1:1255  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro & /opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8066 -r 127.0.0.1:1256  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro & /opt/udp2raw_{{ udp2raw_arch }} -s -l0.0.0.0:8067 -r 127.0.0.1:1257  -k "{{ udp2raw_secret }}" --raw-mode faketcp -a --cipher-mode xor --auth-mode simple --disable-bpf --fix-gro
WorkingDirectory=/opt
RestartSec=10
Restart=always

[Install]
WantedBy=multi-user.target
